# Client/camera-server.service
[Unit]
Description=Camera Server for Thermal Printer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/thermalPrinterRaspy
# Kill any process on port 5001 before starting
ExecStartPre=-/usr/bin/fuser -k 5001/tcp
ExecStartPre=-/bin/sleep 5
ExecStart=/home/pi/thermalPrinterRaspy/venv/bin/python -u /home/pi/thermalPrinterRaspy/007_arducam_qr_system.py --server ws://172.20.10.2:8765
# Also kill port when stopping
ExecStopPost=-/usr/bin/fuser -k 5001/tcp
Restart=always
RestartSec=10
TimeoutStopSec=10
KillMode=mixed
KillSignal=SIGTERM
Environment=DISPLAY=:0
Environment=HOME=/home/pi
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
